# systemd/pihole-group-toggle@.timer
[Unit]
Description=Schedule pihole group toggle (%i)

[Timer]
Persistent=true

[Install]
WantedBy=timers.target
